# rc.local
#!/bin/sh -e

echo 'Server starting...'
sudo python3 /home/pi/Cascade/main.py


exit 0
